# Application_tweaks/initscripts/mailcatcher.upstart
description "mailcatcher"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

chdir /srv

script
    /usr/local/bin/mailcatcher --ip 0.0.0.0
end script
